# packages/atoms-extra/telemetry/app-metrics/node_exporter-bin/node_exporter.service
[Unit]
Description=Node Exporter

[Service]
User=node_exporter
EnvironmentFile=/etc/conf.d/node_exporter
ExecStart=/usr/bin/node_exporter $EXPORTER_ARGS

[Install]
WantedBy=multi-user.target
